# .hive/launchd/ai.csoai.meok-ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.csoai.meok-ui</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>cd /Users/nicholas/clawd/meok/ui &amp;&amp; exec /usr/bin/python3 -m http.server 3000</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/nicholas/clawd/.hive/logs/ai.csoai.meok-ui.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/nicholas/clawd/.hive/logs/ai.csoai.meok-ui.stdout.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>/Users/nicholas/clawd</string>
</dict>
</plist>
