# .hive/launchd/ai.csoai.capital-ascension-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.csoai.capital-ascension-orchestrator</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/nicholas/clawd/sovereign-temple/.venv/bin/python3</string>
		<string>/Users/nicholas/clawd/.hive/scripts/capital_ascension/capital_ascension_orchestrator.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/nicholas/clawd/.hive/logs/ai.csoai.capital-ascension-orchestrator.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/nicholas/clawd/.hive/logs/ai.csoai.capital-ascension-orchestrator.stdout.log</string>
	<key>StartInterval</key>
	<integer>604800</integer>
	<key>WorkingDirectory</key>
	<string>/Users/nicholas/clawd</string>
</dict>
</plist>
